# knowledge-service/systemd/knowledge-backup.timer
[Unit]
Description=Weekly knowledge index backup

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target

