# knowledge-service/systemd/knowledge-tunnel.service
[Unit]
Description=Cloudflare Quick Tunnel for the Knowledge API
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/local/bin/cloudflared tunnel --no-autoupdate --url http://127.0.0.1:8088
Restart=always
RestartSec=5
TimeoutStopSec=15

[Install]
WantedBy=default.target
