# scripts/openclaw-website-session-cleanup.service
[Unit]
Description=Clean expired OpenClaw website sessions

[Service]
Type=oneshot
ExecStart=%h/.local/bin/openclaw-website-session-cleanup
